# app/build/php/symfony/services/nginx/run
#!/usr/bin/env bash

set -eu
/usr/sbin/nginx 
