# app/build/php/symfony/services/php-fpm/run
#!/usr/bin/env bash

set -eu

exec /usr/sbin/php-fpm7.2 -y /etc/php/7.2/fpm/php-fpm.conf 2>&1
